# docker-plugins-xml/os-services/resources/services.conf
[program:sshd]
command=/usr/sbin/sshd -De
startsecs=5
stderr_logfile = /var/log/supervisor/sshd-stderr.log
stdout_logfile = /var/log/supervisor/sshd-stdout.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
startsecs=5
stderr_logfile = /var/log/supervisor/nginx-stderr.log
stdout_logfile = /var/log/supervisor/nginx-stdout.log

[program:app-uwsgi]
; --catch-exceptions for debugging only!
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini --uid uwsgi --gid uwsgi
startsecs=5
stderr_logfile = /var/log/supervisor/app-uwsgi-stderr.log
stdout_logfile = /var/log/supervisor/app-uwsgi-stdout.log
